# selahseedcore/services/selah-wifi-wake.service
[Unit]
Description=SelahSeedCore — BCM4350 PCIe Wake Sequence
After=systemd-modules-load.service
Before=NetworkManager.service
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c '\
  rmmod brcmfmac 2>/dev/null || true; \
  sleep 2; \
  echo 1 > /sys/bus/pci/devices/0000:02:00.0/remove; \
  sleep 3; \
  echo 1 > /sys/bus/pci/rescan; \
  sleep 3; \
  modprobe brcmfmac || true'
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
